# utils/dashcam.service.template
[Unit]
Description=Dashcam
StartLimitBurst=10
StartLimitIntervalSec=0
StartLimitAction=reboot

[Install]
WantedBy=multi-user.target

[Service]
Type=simple
Restart=always
RestartSec=10
